# supervisord.conf
[supervisord]
nodaemon=true

[program:postwycker]
command=/usr/bin/python /opt/postwycker/collector.py
